# ndn_python_repo/ndn-python-repo.service
[Unit]
Description=NDN Repo Service

[Service]
Type=simple
ExecStart=/usr/local/bin/ndn-python-repo

[Install]
WantedBy=multi-user.target
